# arroyo.conf
description "Arroyo"

start on runlevel [2345]
stop on runlevel [!2345]

chdir /srv/arroyo
exec su - chris -c '/usr/bin/python3 -u /srv/arroyo/arroyo.py'
