# liana.conf
description "Liana"

start on runlevel [2345]
stop on runlevel [!2345]

chdir /srv/liana
exec su - chris -c '/usr/bin/python3 -u /srv/liana/liana.py'
